# nginx/rootfs/etc/services.d/fpm/run
#!/usr/bin/env bash
set -e
exec /usr/sbin/php-fpm8 \
    --pid /var/run/php-fpm8/php-fpm8.pid \
    --prefix /var/run/php-fpm8 \
    --fpm-config /etc/php8/php-fpm.conf
